# scripts/com.emma.org-future-insights.bailian.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.emma.org-future-insights.bailian</string>

    <!-- 06:30 跑：抓取已在 06:00 由原 plist 完成 -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-m</string>
        <string>bailian.pipeline</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>OFI_PROJECT_ROOT</key>
        <string>/Users/emmah/org-future-insights</string>
        <key>PYTHONPATH</key>
        <string>/Users/emmah/Desktop/AI学习/Qoder文档/Emmafolder/组织演变</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <!-- 必须 cd 到工作区，bailian 包才能被 import -->
    <key>WorkingDirectory</key>
    <string>/Users/emmah/Desktop/AI学习/Qoder文档/Emmafolder/组织演变</string>

    <key>StandardOutPath</key>
    <string>/Users/emmah/org-future-insights/daily-raw/_logs/launchd.bailian.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/emmah/org-future-insights/daily-raw/_logs/launchd.bailian.err.log</string>

    <key>Nice</key>
    <integer>10</integer>

    <key>ThrottleInterval</key>
    <integer>3600</integer>
</dict>
</plist>
